# systemctl/user/elden-ring-backup.service
[Unit]
Description=Rotating 5 Elden Ring save backups

[Service]
ExecStart=/home/kevon/scripts/elden-ring-backup.sh
